# web/vagrant/supervisor.conf
[program:server]
command=nim c -r server
directory=/var/www
numprocs=1
stdout_logfile=/var/www/daemon.log
autostart=true
autorestart=true
user=vagrant
stopsignal=KILL

